# pd-mapper/DEBIAN/postrm
#!/bin/sh
set -e

case "$1" in
    remove|purge)
        if command -v systemctl >/dev/null 2>&1; then
            systemctl stop pd-mapper.service || true
            systemctl disable pd-mapper.service || true
            systemctl daemon-reload
        fi
    ;;
esac

exit 0